# hubsaas-frontend.service
[Unit]
Description=ngrok tunnel HubSaaS frontend (:3020)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__HOME__

ExecStartPre=__HOME__/debian/wait-ngrok-port.sh 3020
ExecStart=__HOME__/debian/start-ngrok-tunnel.sh frontend
ExecStop=/bin/kill -SIGINT $MAINPID

Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=20

StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
